# debian/em-mysqlbackup.cron.d
18 */6 * * *	root	[ -x /usr/bin/em-mysqlbackup ] && /usr/bin/em-mysqlbackup
